# deploy/supervisord_opac.ini
[supervisord]
nodaemon = true
pidfile = /var/run/supervisord.pid
logfile = /var/log/supervisord.log
loglevel = info

[program:nginx]
user=nginx
startsecs = 5
autostart = true
autorestart = true
loglevel = info
command = /usr/sbin/nginx

[program:opac_gunicorn]
user=nginx
autostart=true
autorestart=true
loglevel = info
stdout_logfile = /var/log/supervisor_opac_gunicorn.out.log
stderr_logfile = /var/log/supervisor_opac_gunicorn.err.log
command = gunicorn --workers 3 --bind unix:/var/run/nginx/opac.sock manager:app --chdir=/var/www/opac --log-level INFO

[inet_http_server]
port=127.0.0.1:9001
